# shell/profile.d/01.dircolors.sh
#!/bin/bash

# Dircolors
eval `dircolors ~/.config/dotfiles/dircolors`

# Midnight Commander
export MC_SKIN=${HOME}/.config/dotfiles/mc.ini
