# shell/profile.d/01.dircolors.sh
#!/bin/bash
eval `dircolors ~/.config/dotfiles/dircolors`
